# supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0766

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
user = root
childlogdir = /var/www/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:php-fpm]
command = /usr/local/sbin/php-fpm
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true

[program:cron]
command = /usr/sbin/cron -f
autorestart = true
user = root
stdout_events_enabled = true
stderr_events_enabled = true

[program:sf-messenger]
command = php /var/www/html/bin/console messenger:consume -l 150 -m 240M -t 2400
autorestart = true
user = www-data
environment = HOME="/var/www"
stdout_events_enabled = true
stderr_events_enabled = true

;[eventlistener:stdout]
;command = supervisor_stdout
;buffer_size = 100
;events = PROCESS_LOG
;result_handler = supervisor_stdout:event_handler

